# crontab/crontab-laptop
# MIN HOUR DAY MONTH DAYOFWEEK  COMMAND
0 * * * * ~/Private/bin/clearFlashCookies.sh
